# urm.service.in
[Unit]
Description=URM Service
After=multi-user.target

[Service]
Restart=on-failure
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/urm
TimeoutStopSec=2
Slice=focused.slice
Delete=yes

[Install]
WantedBy=multi-user.target
