# systune/debian/qcom-systune.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl enable qcom-systune.service || true
fi

exit 0

