# debian/qcom-systune.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl enable resource-tuner.service || true
fi

exit 0

